# rootfs/etc/s6-overlay/s6-rc.d/nodered/finish
#!/bin/bash
# Exit code 256 means the process was killed by a signal
if [ "$1" -eq 256 ]; then
    s6-svscanctl -t /var/run/s6/services
fi
